# scripts/anime-gallery.service
[Unit]
Description=MoeGallery
After=network-online.target mysql.service mysqld.service mariadb.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/moegallery
ExecStart=/opt/moegallery/venv/bin/python /opt/moegallery/scripts/moegallery_launcher.py --app-dir /opt/moegallery --host 127.0.0.1 --port 8111
Restart=always
RestartSec=3
TimeoutStopSec=45
UMask=0027
NoNewPrivileges=true
User=moegallery
Group=moegallery

[Install]
WantedBy=multi-user.target
